# launchd/com.agent-memory.consolidate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-memory.consolidate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/memory</string>
        <string>consolidate</string>
    </array>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>StandardOutPath</key>
    <string>/tmp/agent-memory-consolidate.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agent-memory-consolidate.log</string>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
